# recipes-summit/summit-initdata/files/rootfs-additions-secure/usr/lib/systemd/system/systemd-timedated.service.d/01-rw.conf
[Service]
ReadWritePaths=/etc /data/misc
